# systemd/toggle.service
[Unit]
Description=The Redeem/Replicape user interface

[Service]
Type=simple
ExecStart=/usr/local/bin/toggle

[Install]
WantedBy=basic.target

