# debian/postinst
#!/bin/sh

# Enable redeem.service
systemctl enable redeem.service
