# debian/postinst
#!/bin/sh

# Enable toggle.service
systemctl enable toggle.service
systemctl start toggle.service
